# dot_config/systemd/user/shutdown-on-morning.service.tmpl
{{ if eq .chezmoi.hostname "PROTEUSNOTEBOOK" -}}
[Unit]
Description=Shutdown system

[Service]
ExecStart=/usr/bin/systemctl poweroff
{{- end -}}
