# dot_config/systemd/user/waybar.service.d/override.conf.tmpl
{{ if lookPath "waybar" -}}
[Service]
RestartSec=3sec
Slice=background.slice
{{- end -}}
